# mission-control.service
[Unit]
Description=Mission Control Dashboard for OpenClaw
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/path/to/mission-control
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3333

[Install]
WantedBy=multi-user.target